# files/minecraft-on-demand.socket
[Unit]
Description=On-demand Minecraft socket

[Socket]
ListenStream=25565
Accept=yes

[Install]
WantedBy=sockets.target
